# mailikan.service
[Unit]
Description=Mailikan Email Kanban Application
Documentation=https://github.com/Plobli/Mailikan
After=network.target

[Service]
Type=notify
User=christopher
Group=christopher
WorkingDirectory=/opt/mailikan
Environment=NODE_ENV=production
Environment=PORT=3001
ExecStart=/usr/bin/node server.js
Restart=on-failure
RestartSec=10
TimeoutStartSec=60
TimeoutStopSec=30
StandardOutput=journal
StandardError=journal
SyslogIdentifier=mailikan

# Security Settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/mailikan/data
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

[Install]
WantedBy=multi-user.target
